# Proyecto/AWS_IaC/configuration_files/zookeeper/supervisord_zk_process.conf
[program:zookeeper]
command=/zk/bin/zkServer.sh start-foreground
autorestart=true
stopsignal=KILL
